# pi-gen-pantry/00-shop/files/shop-inventory.service
# TODO: should use environment variables from config instead of literals
[Unit]
Description=Rowan Shop & Pantry Inventory System
After=network.target

[Service]
Type=simple
User=pantry
Group=www-data
RuntimeDirectory=shop-inventory
WorkingDirectory=/var/www/shop-inventory
EnvironmentFile=/etc/shop-inventory/config
ExecStart=/bin/bash start.sh
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Limit resource usage
LimitNOFILE=1048576
LimitNPROC=64

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
PrivateUsers=true

[Install]
WantedBy=multi-user.target
