# stageCustom/00-shop/files/shop-backup.service
[Unit]
Description=Shop Inventory Database Backup Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/shop-backup.sh
User=pantry
Group=www-data
Environment="DJANGO_BACKUP_PASSWORD=${DJANGO_BACKUP_PASSWORD}"

[Install]
WantedBy=multi-user.target
